# debian/naemon-core.prerm
#!/bin/sh -e

if [ -x "/etc/init.d/naemon" ]; then
    if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
        invoke-rc.d naemon stop || true
    else
        /etc/init.d/naemon stop || true
    fi
    # Remove rc.d scripts before uninstalling
    update-rc.d -f naemon remove
fi

#DEBHELPER#
